# BT_tethering_client/tethering.service
[Unit]
Description=Starts tethering service to connect to tatblet

[Service]
ExecStart=/home/xxx/tethering.sh

[Install]
WantedBy=multi-user.target